# templates/etc/profile.d/java.sh
JAVA_HOME="{{ java_install_folder }}/java"
PATH="$JAVA_HOME/bin:$PATH"
